# zipit-base-files/files/etc/init.d/local
#!/bin/sh /etc/rc.common
#
# Process user commands
#

START=99

start() {
	[ -f /etc/rc.local ] && /bin/sh /etc/rc.local
}
