# zipit-base-files/files/etc/inittab
# Rescan devices
#::sysinit:/sbin/mdev -s

# Set the hostname
::sysinit:/bin/hostname -F /etc/hostname

# Set Zipit Z2 Console keymap and font
::sysinit:/usr/sbin/loadfont < /etc/5x8iz2s.psf
::sysinit:/sbin/loadkmap < /etc/z2.kmap

# Run any rc scripts
::sysinit:/etc/init.d/rcS S boot

# Reduce kernel messages on the console
::sysinit:/bin/echo "4 1 1 7" > /proc/sys/kernel/printk

# Stuff to do before rebooting
::shutdown:/etc/init.d/rcS K shutdown
::shutdown:/sbin/swapoff -a
::shutdown:/bin/sync
::shutdown:/bin/umount -a -r

# Start the tty's
ttyS2::askfirst:/bin/ash --login
tty1::respawn:/usr/bin/z2-shell
#tty1::respawn:/bin/ash --login
tty2::respawn:/usr/sbin/ebindkeys -d -f /etc/ebindkeysrc
